# flatpack2.service
[Unit]
Description=flatpack2 - Eltek Flatpack2 CAN controller
Documentation=https://github.com/mti-sk/flatpack2
After=network.target

[Service]
# Matches the layout created by install.sh (/opt/flatpack2 + venv,
# config in /etc/flatpack2). Adjust paths for a custom installation.
Type=simple
ExecStart=/opt/flatpack2/venv/bin/python /opt/flatpack2/flatpack2.py --config /etc/flatpack2/flatpack2.conf
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=flatpack2

# Run as specific user (change as needed)
# User=flatpack2
# Group=dialout

# Give access to serial ports
SupplementaryGroups=dialout

# Working directory (relative log file path resolves here)
WorkingDirectory=/opt/flatpack2

# Security hardening
# NOTE: PrivateTmp must stay disabled - the PTY symlink /tmp/flatpack2.pty
# would otherwise be created in a private namespace and 'screen' could not
# attach to it from a user session.
ProtectSystem=full
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
